# phingr-device/rpi/setup/phingr-capture.service
[Unit]
Description=phingr Capture Server (Pi Camera)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/phingr/capture_server.py --width 1280 --height 720 --fps 15
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
